# packages/secubox-traffic/debian/secubox-traffic.service
[Unit]
Description=SecuBox Traffic Shaper API
After=network.target secubox-core.service
Wants=secubox-core.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/usr/lib/secubox/traffic
ExecStart=/usr/bin/python3 -m uvicorn api.main:app --uds /run/secubox/traffic.sock --log-level warning
Restart=on-failure
RestartSec=5
PrivateTmp=true

[Install]
WantedBy=multi-user.target
