# remote-ui/round/secubox-eye-agent.service
[Unit]
Description=SecuBox Eye Remote Agent
After=network.target pigpiod.service
Wants=pigpiod.service

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /usr/lib/secubox-eye/agent/main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Runtime directory
RuntimeDirectory=secubox-eye
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
