# packages/secubox-wireguard/debian/prerm
#!/bin/sh
set -e

case "$1" in
    remove|upgrade|deconfigure)
        systemctl stop secubox-wireguard.service || true
        systemctl disable secubox-wireguard.service || true
        ;;
esac

#DEBHELPER#

exit 0
