# remote-ui/round/files/etc/systemd/system/dnsmasq.service.d/secubox-eye.conf
[Unit]
After=secubox-eye-gadget.service
Wants=secubox-eye-gadget.service

[Service]
# Load environment for optional DHCP toggle
EnvironmentFile=-/etc/secubox/eye-remote/tftp.env
