# packages/secubox-vhost/debian/prerm
#!/bin/bash
set -e
case "$1" in
  remove|upgrade)
    rm -f /etc/nginx/secubox.d/vhost.conf
    systemctl stop    secubox-vhost.service 2>/dev/null || true
    systemctl disable secubox-vhost.service 2>/dev/null || true
    systemctl reload nginx 2>/dev/null || true
    ;;
esac
#DEBHELPER#
