# packages/secubox-vhost/debian/prerm
#!/bin/bash
set -e
case "$1" in
  remove|upgrade)
    systemctl stop    secubox-vhost.service 2>/dev/null || true
    systemctl disable secubox-vhost.service 2>/dev/null || true
    ;;
esac
#DEBHELPER#
