# packages/secubox-ksm/debian/prerm
#!/bin/bash
set -e
case "$1" in
  remove|deconfigure)
    systemctl stop secubox-ksm.service || true
    systemctl disable secubox-ksm.service || true
    ;;
esac
#DEBHELPER#
