# packages/secubox-rezapp/debian/prerm
#!/bin/sh
set -e
if [ "$1" = "remove" ]; then
    systemctl stop secubox-rezapp.service || true
    systemctl disable secubox-rezapp.service || true
fi
#DEBHELPER#
exit 0
