# packages/secubox-mail/debian/secubox-mail.service
[Unit]
Description=SecuBox Mail API
After=network.target secubox-core.service
Requires=secubox-core.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/usr/lib/secubox/mail
ExecStart=/usr/bin/python3 -m uvicorn api.main:app --uds /run/secubox/mail.sock --log-level warning
Restart=on-failure
RestartSec=5
UMask=0000

[Install]
WantedBy=multi-user.target
