# packages/secubox-waf/debian/prerm
#!/bin/sh
set -e
if [ "$1" = "remove" ]; then
    systemctl stop secubox-waf.service || true
    systemctl disable secubox-waf.service || true
    systemctl reload nginx 2>/dev/null || true
fi
#DEBHELPER#
exit 0
