# packages/secubox-voip/debian/postinst
#!/bin/bash
set -e
case "$1" in
  configure)
    install -d -o secubox -g secubox -m 750 /var/lib/secubox/voip
    systemctl daemon-reload
    systemctl enable secubox-voip.service
    systemctl start secubox-voip.service || true
    systemctl reload nginx 2>/dev/null || true
    ;;
esac
#DEBHELPER#
