# packages/secubox-newsbin/debian/secubox-newsbin.service
[Unit]
Description=SecuBox Newsbin API
After=network.target secubox-core.service

[Service]
Type=simple
User=root
WorkingDirectory=/usr/lib/secubox/newsbin
ExecStart=/usr/bin/uvicorn api.main:app --uds /run/secubox/newsbin.sock --log-level warning
UMask=0117
Restart=on-failure
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/run/secubox /var/lib/secubox /var/lib/lxc

[Install]
WantedBy=multi-user.target
