# packages/secubox-newsbin/debian/secubox-newsbin.service
[Unit]
Description=SecuBox Newsbin API
After=network.target secubox-core.service secubox-runtime.service docker.service
Requires=secubox-runtime.service
Wants=docker.service
# Only start if explicitly enabled (SABnzbd installed)
ConditionPathExists=/etc/secubox/newsbin/enabled

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/usr/lib/secubox/newsbin
ExecStart=/usr/bin/python3 -m uvicorn api.main:app --uds /run/secubox/newsbin.sock --log-level warning
ExecStartPost=-/bin/sleep 1
ExecStartPost=-/bin/chmod 660 /run/secubox/newsbin.sock
ExecStartPost=-/bin/chown secubox:secubox /run/secubox/newsbin.sock
Restart=on-failure
RestartSec=5
PrivateTmp=true

[Install]
WantedBy=multi-user.target
