# packages/secubox-webradio/debian/secubox-webradio.service
[Unit]
Description=SecuBox WebRadio API
After=network.target secubox-core.service

[Service]
Type=simple
User=root
WorkingDirectory=/usr/lib/secubox/webradio
ExecStart=/usr/bin/uvicorn api.main:app --uds /run/secubox/webradio.sock --log-level warning
UMask=0117
Restart=on-failure
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/run/secubox /var/lib/secubox

[Install]
WantedBy=multi-user.target
