# packages/secubox-publish/debian/postinst
#!/bin/sh
set -e
if [ "$1" = "configure" ]; then
    systemctl daemon-reload
    systemctl enable secubox-publish.service || true
    systemctl start secubox-publish.service || true
fi
#DEBHELPER#
exit 0
