# packages/secubox-nettweak/debian/prerm
#!/bin/bash
set -e
case "$1" in
  remove|deconfigure)
    systemctl stop secubox-nettweak.service || true
    systemctl disable secubox-nettweak.service || true
    ;;
esac
#DEBHELPER#
