# packages/secubox-dpi/debian/secubox-dpi/usr/lib/systemd/system/secubox-dpi.service
[Unit]
Description=SecuBox DPI Dashboard API
After=network.target secubox-core.service netifyd.service
Requires=secubox-core.service
Wants=netifyd.service

[Service]
UMask=0000
Type=simple
User=secubox
Group=secubox
WorkingDirectory=/usr/lib/secubox/dpi
ExecStart=/usr/bin/python3 -m uvicorn api.main:app \
    --uds /run/secubox/dpi.sock \
    --log-level warning
Restart=on-failure
RestartSec=5
PrivateTmp=true

[Install]
WantedBy=multi-user.target
