# packages/secubox-cookies/systemd/secubox-cookies.service
[Unit]
Description=SecuBox Cookie Tracker API
After=network.target secubox-core.service
Wants=secubox-core.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/usr/lib/secubox/cookies
ExecStart=/usr/bin/python3 -m uvicorn api.main:app --uds /run/secubox/cookies.sock --log-level warning
Restart=on-failure
RestartSec=5
PrivateTmp=true
RuntimeDirectory=secubox
RuntimeDirectoryPreserve=yes

[Install]
WantedBy=multi-user.target
