# packages/secubox-mcp-server/debian/postinst
#!/bin/sh
set -e
if [ "$1" = "configure" ]; then
    systemctl daemon-reload
    systemctl enable secubox-mcp-server.service || true
    systemctl start secubox-mcp-server.service || true
fi
#DEBHELPER#
exit 0
