# packages/secubox-rtty/debian/prerm
#!/bin/sh
set -e

case "$1" in
    remove|upgrade|deconfigure)
        systemctl stop secubox-rtty.service 2>/dev/null || true
        systemctl disable secubox-rtty.service 2>/dev/null || true
        ;;
esac

#DEBHELPER#

exit 0
