# packages/secubox-redroid/debian/secubox-redroid.service
[Unit]
# Only start if explicitly enabled (backend installed)
ConditionPathExists=/etc/secubox/redroid/enabled
Description=SecuBox Redroid API
After=network.target secubox-core.service

[Service]
Type=simple
User=root
WorkingDirectory=/usr/lib/secubox/redroid
ExecStart=/usr/local/bin/uvicorn api.main:app --uds /run/secubox/redroid.sock --log-level warning
UMask=0117
Restart=on-failure

[Install]
WantedBy=multi-user.target
