# packages/secubox-magicmirror/debian/secubox-magicmirror.service
[Unit]
# Only start if explicitly enabled (backend installed)
ConditionPathExists=/etc/secubox/magicmirror/enabled
Description=SecuBox MagicMirror API
After=network.target secubox-core.service

[Service]
Type=simple
User=root
WorkingDirectory=/usr/lib/secubox/magicmirror
ExecStart=/usr/bin/uvicorn api.main:app --uds /run/secubox/magicmirror.sock --log-level warning
UMask=0117
Restart=on-failure

[Install]
WantedBy=multi-user.target
