# packages/secubox-droplet/packages/secubox-vhost/debian/prerm
#!/bin/bash
set -e
case "$1" in
  remove|upgrade|deconfigure)
    systemctl stop secubox-vhost.service || true
    systemctl disable secubox-vhost.service || true
    ;;
esac
#DEBHELPER#
