# packages/secubox-reporter/debian/prerm
#!/bin/bash
set -e
case "$1" in
  remove|upgrade|deconfigure)
    systemctl stop secubox-reporter.service || true
    systemctl disable secubox-reporter.service || true
    ;;
esac
#DEBHELPER#
