# packages/secubox-backup/debian/prerm
#!/bin/sh
set -e
if [ "$1" = "remove" ]; then
    systemctl stop secubox-backup.service || true
    systemctl disable secubox-backup.service || true
fi
#DEBHELPER#
exit 0
