# packages/secubox-netdata/debian/prerm
#!/bin/bash
set -e
case "$1" in
  remove|upgrade)
    systemctl stop    secubox-netdata.service 2>/dev/null || true
    systemctl disable secubox-netdata.service 2>/dev/null || true
    ;;
esac
#DEBHELPER#
