# packages/secubox-hexo/debian/postinst
#!/bin/bash
set -e
case "$1" in
  configure)
    install -d -o secubox -g secubox -m 750 /var/lib/secubox/hexo/sites
    install -d -m 755 /var/www/hexo
    systemctl daemon-reload
    systemctl enable secubox-hexo.service
    systemctl start secubox-hexo.service || true
    systemctl reload nginx 2>/dev/null || true
    ;;
esac
#DEBHELPER#
