# packages/secubox-daemon/debian/secubox-daemon.prerm
#!/bin/bash
set -e

case "$1" in
  remove|purge)
    systemctl stop secuboxd.service || true
    systemctl disable secuboxd.service || true
    ;;
esac
#DEBHELPER#
