# packages/secubox-ui-manager/debian/secubox-ui-health.timer
[Unit]
Description=SecuBox UI Health Check Timer
Documentation=https://secubox.in/docs/ui-manager

[Timer]
OnBootSec=30
OnUnitActiveSec=30
Unit=secubox-ui-health.service

[Install]
WantedBy=timers.target
