# packages/secubox-vault/debian/secubox-vault.service
[Unit]
Description=SecuBox Vault API
After=network.target secubox-core.service

[Service]
Type=simple
User=root
WorkingDirectory=/usr/lib/secubox/vault
ExecStart=/usr/local/bin/uvicorn api.main:app --uds /run/secubox/vault.sock --log-level warning
UMask=0117
Restart=on-failure

[Install]
WantedBy=multi-user.target
