# packages/secubox-zigbee/debian/secubox-zigbee.service
[Unit]
Description=SecuBox Zigbee2MQTT Gateway API
After=network.target secubox-runtime.service
Requires=secubox-runtime.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/usr/share/secubox/zigbee
ExecStart=/usr/bin/uvicorn api.main:app --uds /run/secubox/zigbee.sock --workers 1
ExecStartPost=-/bin/sleep 1
ExecStartPost=-/bin/chmod 660 /run/secubox/zigbee.sock
ExecStartPost=-/bin/chown secubox:secubox /run/secubox/zigbee.sock
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=no
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/run/secubox /var/lib/secubox /etc/secubox /srv/zigbee2mqtt

[Install]
WantedBy=multi-user.target
