# packages/secubox-peertube/debian/prerm
#!/bin/bash
set -e
case "$1" in
  remove|upgrade)
    # Stop and disable service
    systemctl stop    secubox-peertube.service 2>/dev/null || true
    systemctl disable secubox-peertube.service 2>/dev/null || true
    # Reload nginx to remove location
    systemctl reload nginx 2>/dev/null || true
    ;;
esac
#DEBHELPER#
