# packages/secubox-peertube/debian/secubox-peertube.service
[Unit]
Description=SecuBox PeerTube API
After=network.target secubox-core.service

[Service]
Type=simple
User=root
WorkingDirectory=/usr/lib/secubox/peertube
ExecStart=/usr/bin/uvicorn api.main:app --uds /run/secubox/peertube.sock --log-level warning
UMask=0117
Restart=on-failure
PrivateTmp=true
RuntimeDirectory=secubox
RuntimeDirectoryMode=0775
ProtectSystem=full
ReadWritePaths=/run/secubox /var/lib/secubox /var/lib/lxc

[Install]
WantedBy=multi-user.target
