# packages/secubox-crowdsec/debian/prerm
#!/bin/bash
set -e
case "$1" in
  remove|upgrade)
    rm -f /etc/nginx/secubox.d/crowdsec.conf
    systemctl stop    secubox-crowdsec.service 2>/dev/null || true
    systemctl disable secubox-crowdsec.service 2>/dev/null || true
    systemctl reload nginx 2>/dev/null || true
    ;;
esac
#DEBHELPER#
