# packages/secubox-smtp-relay/debian/prerm
#!/bin/bash
set -e
case "$1" in
  remove|upgrade|deconfigure)
    systemctl stop secubox-smtp-relay.service || true
    systemctl disable secubox-smtp-relay.service || true
    ;;
esac
#DEBHELPER#
