# packages/secubox-console/debian/prerm
#!/bin/bash
set -e

case "$1" in
    remove|upgrade|deconfigure)
        # Stop console service if running
        systemctl stop secubox-console.service 2>/dev/null || true
        systemctl disable secubox-console.service 2>/dev/null || true
        ;;
esac

#DEBHELPER#

exit 0
