# packages/secubox-console/debian/secubox-console.service
[Unit]
Description=SecuBox Console TUI
Documentation=https://cybermind.fr/secubox
After=multi-user.target
ConditionPathExists=/usr/bin/secubox-console

[Service]
Type=simple
ExecStart=/usr/bin/secubox-console
StandardInput=tty
StandardOutput=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
