# packages/secubox-netdiag/debian/prerm
#!/bin/bash
set -e
case "$1" in
  remove|purge)
    systemctl stop secubox-netdiag.service 2>/dev/null || true
    systemctl disable secubox-netdiag.service 2>/dev/null || true
    ;;
esac
#DEBHELPER#
