# packages/secubox-webmail/debian/prerm
#!/bin/sh
set -e
if [ "$1" = "remove" ]; then
    systemctl stop secubox-webmail.service || true
    systemctl disable secubox-webmail.service || true
    systemctl reload nginx 2>/dev/null || true
fi
#DEBHELPER#
exit 0
