# packages/secubox-metabolizer/debian/prerm
#!/bin/bash
set -e
case "$1" in
  remove|upgrade|deconfigure)
    systemctl stop secubox-metabolizer.service || true
    systemctl disable secubox-metabolizer.service || true
    ;;
esac
#DEBHELPER#
