# packages/secubox-threat-analyst/debian/prerm
#!/bin/sh
set -e
if [ "$1" = "remove" ]; then
    systemctl stop secubox-threat-analyst.service || true
    systemctl disable secubox-threat-analyst.service || true
fi
#DEBHELPER#
exit 0
