# packages/secubox-iot-guard/debian/prerm
#!/bin/sh
set -e
if [ "$1" = "remove" ]; then
    systemctl stop secubox-iot-guard.service || true
    systemctl disable secubox-iot-guard.service || true
fi
#DEBHELPER#
exit 0
