# packages/secubox-watchdog/debian/prerm
#!/bin/sh
set -e
if [ "$1" = "remove" ]; then
    systemctl stop secubox-watchdog.service || true
    systemctl disable secubox-watchdog.service || true
fi
#DEBHELPER#
exit 0
