# image/systemd/secubox-kiosk-wayland.service
[Unit]
Description=SecuBox Kiosk Mode (Wayland/Cage)
Documentation=https://secubox.in/docs/kiosk

# Proper dependency chain - wait for everything to be ready
After=systemd-user-sessions.service
After=plymouth-quit-wait.service
After=network-online.target
After=nginx.service
After=systemd-networkd.service
After=secubox-firstboot.service

# Soft dependencies (don't fail if missing)
Wants=network-online.target
Wants=nginx.service

# Conflicts with console on tty1
Conflicts=getty@tty1.service

# Only start if graphical target is active
PartOf=graphical.target

[Service]
Type=simple

# Check if kiosk is enabled (run as root with + prefix)
ExecStartPre=+/bin/sh -c '[ -f /var/lib/secubox/.kiosk-enabled ] || exit 1'

# Use the robust launcher script (runs as root, handles user switching internally)
ExecStart=/usr/sbin/secubox-kiosk-launcher wayland 7

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Generous timeout for slow hardware/VMs
TimeoutStartSec=180
TimeoutStopSec=30

# Resource limits
MemoryMax=1G
TasksMax=200

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=secubox-kiosk

[Install]
WantedBy=graphical.target
