# image/systemd/secubox-kiosk.service
[Unit]
Description=SecuBox Kiosk Mode (X11/Chromium)
Documentation=https://secubox.in/docs/kiosk

# Start after multi-user (not graphical - we ARE the graphical part)
After=multi-user.target
After=systemd-user-sessions.service

# CRITICAL: Wait for X11 driver auto-detection (VirtualBox VMSVGA, etc.)
Wants=secubox-x11-setup.service
After=secubox-x11-setup.service

# Soft dependency on nginx (works without it)
Wants=nginx.service
After=nginx.service

# Take over VT7
Conflicts=getty@tty7.service

# Rate-limit restarts to prevent VT freeze from rapid cycling
StartLimitBurst=3
StartLimitIntervalSec=120

[Service]
Type=simple

# Clean VT state before starting
ExecStartPre=+/usr/bin/chvt 7
ExecStartPre=+/usr/bin/clear

# Check if kiosk is enabled
ExecStartPre=+/bin/sh -c '[ -f /var/lib/secubox/.kiosk-enabled ] || exit 1'

# Run the launcher (execs to xinit, keeps running)
ExecStart=/usr/sbin/secubox-kiosk-launcher

# Clean up X on stop
ExecStopPost=-/usr/bin/pkill -9 Xorg
ExecStopPost=-/usr/bin/chvt 1

# Use VT7 for X11
TTYPath=/dev/tty7
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes

# Restart on failure but with longer delay
Restart=on-failure
RestartSec=15

# Give X time to start
TimeoutStartSec=120
TimeoutStopSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=secubox-kiosk

[Install]
# Install into multi-user.target, not graphical.target
WantedBy=multi-user.target
