# packages/secubox-config-advisor/debian/prerm
#!/bin/sh
set -e
if [ "$1" = "remove" ]; then
    systemctl stop secubox-config-advisor.service || true
    systemctl disable secubox-config-advisor.service || true
fi
#DEBHELPER#
exit 0
