# packages/secubox-jabber/debian/secubox-jabber.service
[Unit]
# Only start if explicitly enabled (backend installed)
ConditionPathExists=/etc/secubox/jabber/enabled
Description=SecuBox Jabber API
After=network.target secubox-core.service

[Service]
Type=simple
User=root
WorkingDirectory=/usr/lib/secubox/jabber
ExecStart=/usr/bin/uvicorn api.main:app --uds /run/secubox/jabber.sock --log-level warning
UMask=0117
Restart=on-failure

[Install]
WantedBy=multi-user.target
