# remote-ui/round/files/etc/systemd/system/usb-network.service
[Unit]
Description=Configure USB OTG Network
After=secubox-eye-gadget.service
Requires=secubox-eye-gadget.service

[Service]
Type=oneshot
RemainAfterExit=yes
# Wait for gadget to fully initialize
ExecStartPre=/bin/sleep 3
ExecStart=/usr/local/bin/usb-network-up.sh
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
