# packages/secubox-meshname/debian/prerm
#!/bin/sh
set -e
if [ "$1" = "remove" ]; then
    systemctl stop secubox-meshname.service || true
    systemctl disable secubox-meshname.service || true
fi
#DEBHELPER#
exit 0
