# packages/secubox-roadmap/debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        systemctl daemon-reload
        systemctl enable secubox-roadmap.service || true
        systemctl restart secubox-roadmap.service || true
        systemctl reload nginx 2>/dev/null || true
        ;;
esac

#DEBHELPER#

exit 0
