# packages/secubox-mediaflow/debian/prerm
#!/bin/bash
set -e
case "$1" in
  remove|upgrade)
    systemctl stop    secubox-mediaflow.service 2>/dev/null || true
    systemctl disable secubox-mediaflow.service 2>/dev/null || true
    ;;
esac
#DEBHELPER#
