# packages/secubox-mediaflow/debian/prerm
#!/bin/bash
set -e
case "$1" in
  remove|deconfigure)
    rm -f /etc/nginx/secubox.d/mediaflow.conf
    systemctl reload nginx 2>/dev/null || true
    ;;
  upgrade)
    # Don't remove nginx config on upgrade
    ;;
esac
systemctl stop    secubox-mediaflow.service 2>/dev/null || true
systemctl disable secubox-mediaflow.service 2>/dev/null || true
#DEBHELPER#
