# packages/secubox-interceptor/debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Create config directory
        mkdir -p /etc/secubox
        mkdir -p /etc/secubox/interceptor

        # Create run directory
        mkdir -p /run/secubox

        # Create data directories
        mkdir -p /var/lib/secubox/interceptor
        mkdir -p /var/lib/secubox/interceptor/flows
        mkdir -p /var/lib/secubox/interceptor/recordings
        mkdir -p /var/lib/secubox/interceptor/replays
        mkdir -p /tmp/secubox

        # Create log directory
        mkdir -p /var/log/secubox

        # Generate CA certificate if not exists
        if [ ! -f /etc/secubox/interceptor/ca.crt ]; then
            openssl req -x509 -newkey rsa:4096 -nodes \
                -keyout /etc/secubox/interceptor/ca.key \
                -out /etc/secubox/interceptor/ca.crt \
                -days 3650 \
                -subj "/CN=SecuBox Interceptor CA/O=SecuBox/C=FR" 2>/dev/null || true
            chmod 600 /etc/secubox/interceptor/ca.key
        fi

        # Enable and start service
        systemctl daemon-reload
        systemctl enable secubox-interceptor.service
        systemctl start secubox-interceptor.service || true
        ;;
esac

#DEBHELPER#

exit 0
