# packages/secubox-portal/debian/prerm
#!/bin/bash
set -e
case "$1" in
  remove|upgrade)
    rm -f /etc/nginx/secubox.d/portal.conf
    systemctl stop    secubox-portal.service 2>/dev/null || true
    systemctl disable secubox-portal.service 2>/dev/null || true
    systemctl reload nginx 2>/dev/null || true
    ;;
esac
#DEBHELPER#
