# packages/secubox-network-anomaly/debian/prerm
#!/bin/sh
set -e
if [ "$1" = "remove" ]; then
    systemctl stop secubox-network-anomaly.service || true
    systemctl disable secubox-network-anomaly.service || true
fi
#DEBHELPER#
exit 0
