# packages/secubox-hub/debian/prerm
#!/bin/bash
set -e
case "$1" in
  remove|upgrade)
    systemctl stop    secubox-hub.service 2>/dev/null || true
    systemctl disable secubox-hub.service 2>/dev/null || true
    ;;
esac
#DEBHELPER#
