# packages/secubox-avatar/debian/prerm
#!/bin/bash
set -e
case "$1" in
  remove|deconfigure)
    systemctl stop secubox-avatar.service || true
    systemctl disable secubox-avatar.service || true
    ;;
esac
#DEBHELPER#
