# packages/secubox-wazuh/debian/secubox-wazuh.service
[Unit]
Description=SecuBox Wazuh API
After=network.target secubox-core.service

[Service]
Type=simple
User=root
WorkingDirectory=/usr/lib/secubox/wazuh
ExecStart=/usr/local/bin/uvicorn api.main:app --uds /run/secubox/wazuh.sock --log-level warning
UMask=0117
Restart=on-failure

[Install]
WantedBy=multi-user.target
