# packages/secubox-exposure/debian/prerm
#!/bin/sh
set -e

case "$1" in
    remove|upgrade|deconfigure)
        systemctl stop secubox-exposure.service || true
        systemctl disable secubox-exposure.service || true
        ;;
esac

#DEBHELPER#

exit 0
