# packages/secubox-exposure/debian/secubox-exposure.service
[Unit]
Description=SecuBox Exposure Manager API
After=network.target secubox-core.service
Wants=secubox-core.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/usr/lib/secubox/exposure
ExecStart=/usr/bin/python3 -m uvicorn api.main:app --uds /run/secubox/exposure.sock --log-level warning
Restart=on-failure
RestartSec=5
PrivateTmp=true

[Install]
WantedBy=multi-user.target
