# image/secubox-preseed.service
[Unit]
Description=SecuBox Preseed Configuration Apply
After=network-pre.target
Before=network.target
ConditionPathExists=/usr/share/secubox/preseed.tar.gz
ConditionPathExists=!/var/lib/secubox/.preseed-applied

[Service]
Type=oneshot
ExecStart=/usr/lib/secubox/preseed-apply.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
